# config/launchd/com.macroclaw.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.macroclaw.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>macroclaw watch</string>
    </array>

    <!-- Start on user login and keep running -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <!-- Restart if the process exits unexpectedly -->
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Wait 30 seconds before restarting after a crash -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Environment variables available to the job -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/ben</string>
        <key>MACROCLAW_DB_PATH</key>
        <string>/Users/ben/projects/macroclaw/data/macroclaw.duckdb</string>
    </dict>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/tmp/macroclaw-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/macroclaw-watcher.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/ben/projects/macroclaw</string>

</dict>
</plist>
